# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/run/supervisord.pid

[program:xvfb]
command=/usr/bin/Xvfb :0 -screen 0 1280x720x24 -ac +extension GLX +render -noreset
priority=10
autorestart=true
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb-error.log

[program:fluxbox]
command=/bin/bash -lc "sleep 2; exec /usr/bin/fluxbox"
environment=DISPLAY=":0",HOME="/home/ubuntu",USER="ubuntu"
user=ubuntu
priority=20
autorestart=true
stdout_logfile=/var/log/supervisor/fluxbox.log
stderr_logfile=/var/log/supervisor/fluxbox-error.log

[program:xterm]
command=/bin/bash -lc "sleep 3; exec /usr/bin/xterm -geometry 110x32+60+60 -fa Monospace -fs 11 -title 'Ubuntu Terminal' -e /bin/bash --login"
environment=DISPLAY=":0",HOME="/home/ubuntu",USER="ubuntu"
user=ubuntu
priority=30
autorestart=true
stdout_logfile=/var/log/supervisor/xterm.log
stderr_logfile=/var/log/supervisor/xterm-error.log

[program:x11vnc]
command=/bin/bash -lc "sleep 3; exec /usr/bin/x11vnc -display :0 -forever -shared -nopw -rfbport 5900 -listen 127.0.0.1 -xkb"
priority=40
autorestart=true
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc-error.log

[program:websockify]
command=/bin/bash -lc "sleep 4; exec /usr/bin/websockify 127.0.0.1:6080 127.0.0.1:5900"
priority=50
autorestart=true
stdout_logfile=/var/log/supervisor/websockify.log
stderr_logfile=/var/log/supervisor/websockify-error.log

[program:ttyd]
command=/usr/bin/ttyd --writable --base-path /terminal-ws --port 7681 --interface 127.0.0.1 --credential ubuntu:ubuntu /usr/bin/su - ubuntu
priority=60
autorestart=true
stdout_logfile=/var/log/supervisor/ttyd.log
stderr_logfile=/var/log/supervisor/ttyd-error.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=100
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx-error.log
