# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/run/supervisord.pid

[program:xvfb]
command=/usr/bin/Xvfb :0 -screen 0 1366x768x24 -ac +extension GLX +render -noreset
priority=10
autorestart=true
startsecs=2
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb-error.log

[program:xfce]
command=/bin/bash -lc "sleep 3; exec dbus-run-session -- startxfce4"
environment=DISPLAY=":0",HOME="/home/user",USER="user",XDG_RUNTIME_DIR="/tmp/runtime-user"
user=user
priority=20
autorestart=true
startsecs=5
stdout_logfile=/var/log/supervisor/xfce.log
stderr_logfile=/var/log/supervisor/xfce-error.log

[program:x11vnc]
command=/bin/bash -lc "sleep 6; exec /usr/bin/x11vnc -display :0 -forever -shared -nopw -rfbport 5900 -listen 127.0.0.1 -xkb -noxdamage"
priority=30
autorestart=true
startsecs=3
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc-error.log

[program:websockify]
command=/bin/bash -lc "sleep 7; exec /usr/bin/websockify 127.0.0.1:6080 127.0.0.1:5900"
priority=40
autorestart=true
startsecs=3
stdout_logfile=/var/log/supervisor/websockify.log
stderr_logfile=/var/log/supervisor/websockify-error.log

[program:ttyd]
command=/usr/bin/ttyd --writable --base-path /terminal-ws --port 7681 --interface 127.0.0.1 --credential user:1234 /usr/bin/su - user
priority=50
autorestart=true
startsecs=3
stdout_logfile=/var/log/supervisor/ttyd.log
stderr_logfile=/var/log/supervisor/ttyd-error.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=100
autorestart=true
startsecs=3
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx-error.log
